# root/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-yt-dlp-add-package/up
/etc/s6-overlay/s6-rc.d/init-mod-jellyfin-yt-dlp-add-package/run
